# puppy-rootfs-template/usr/lib/puppy/etc/xdg/autostart/layered-fs-bookmark.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
NoDisplay=true
Name=layered filesystem bookmark helper
Exec=layered-fs-bookmark
