# puppy-rootfs-template/usr/lib/puppy/etc/xdg/autostart/notify-pupstate.desktop
[Desktop Entry]
Type=Application
Name=Notify PUPSTATE
Exec=notify-pupstate
