# puppy-rootfs-template/usr/lib/puppy/etc/xdg/autostart/pupmode-13.desktop
[Desktop Entry]
Type=Application
Name=PUPMODE 13 Tray
Comment=Tray for write caching mode
Icon=media-floppy
Exec=pm13 trayicon
