# puppy-rootfs-template/usr/lib/puppy/etc/config-template/profile.d/sudo-askpass.sh
export SUDO_ASKPASS=/usr/lib/puppy/bin/puppy-sudo
