# puppy-rootfs-template/usr/lib/puppy/etc/config-template/cron.d/puppy-pupmode13
*/10 * * * * root /usr/lib/puppy/sbin/pupmode13-cron-helper
